# resources/debian/postrm
#!/bin/sh
# postrm script for natstest

set -e

update-rc.d -f natstest remove

#DEBHELPER#

exit 0
